# ansible/roles/base/templates/podman-prune.timer
[Unit]
Description=Run podman prune daily

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
